# openclaw_task_fixed.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-04-18T15:25:00</Date>
    <Author>DESKTOP-EJV14VH\Administrator</Author>
    <URI>\OpenClaw Gateway</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-21-1200399016-3042025109-561304378-500</UserId>
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>Highest</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
    <SettingsEnabled>true</SettingsEnabled>
  </Settings>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>D:\software\nodejs\node.exe</Command>
      <Arguments>D:\appdata\nodejs\global\node_modules\openclaw\dist\index.js gateway --port 18789</Arguments>
      <WorkingDirectory>D:\appdata\nodejs\global\node_modules\openclaw</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
